# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisord.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:postgresql]
command=/var/www/docker/init-postgresql.sh
autostart=true
autorestart=true
priority=1

[program:symfony]
command=symfony server:start --port=8000 --no-tls --allow-all-ip
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/log/symfony.log
stderr_logfile=/var/log/symfony.err
priority=10

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=5